# src/slimpris2.service.in
[Unit]
Description=Squeezebox MPRIS remote control

[Service]
Type=dbus
ExecStart=@bindir@/slimpris2
BusName=org.mpris.MediaPlayer2.squeezebox
Restart=on-failure

[Install]
WantedBy=default.target
