# src/slimpris2.service.in
[Unit]
Description=Squeezebox MPRIS remote control
Documentation=man:slimpris2.service(7)

[Service]
Type=dbus
ExecStart=@bindir@/slimpris2
BusName=org.mpris.MediaPlayer2.slimpris2
Restart=on-failure
ProtectSystem=full

[Install]
WantedBy=default.target
